# router/crontab
* * * * * python3 /app/manage.py runjob check_site_status >> /var/log/cron.log 2>&1
